# conf/etc/services.d/03-logs/run
#!/bin/sh
NAME="LOGS" 
COMMAND="tail -f /var/log/frog/*" 
WAIT="/tmp/logs"

exec service-handler "$NAME" "$COMMAND" 
